# flask/supervisor_db_monitor.conf
[program:db_monitor]
command=python flask/db_monitor.py
environment=PATH="/srv/pytune/venv/pytune/bin"
directory=/srv/pytune
user=sfvue
autostart=true
autorestart=true
priority=991
stopsignal=HUP
